# fs/etc/crontab
@reboot /home/satoshi/portainer_start.sh
